# observatory.wikivisor/wikivisor/mw-sitemap.service
[Unit]
Description=MediaWiki Sitemap Generator

[Service]
ExecStart=/usr/local/bin/mw-sitemap.generator
Nice=10
ProtectSystem=full
User=www-data
OOMScoreAdjust=200
StandardOutput=journal

[Install]
WantedBy=multi-user.target
